# scripts/watson-dms.service
[Unit]
Description=Watson DMS-SGP02 GPS/inertial logger (switch-gated)
# Wait for the filesystem; the unit + GPIO are local, no network needed.
After=local-fs.target

[Service]
Type=simple
User=__USER__
# dialout = serial port, i2c = Multi-IO HAT, gpio = optional direct-GPIO source.
SupplementaryGroups=dialout i2c gpio
WorkingDirectory=__APPDIR__
# Switch-gated collection: process stays up; a dry-contact switch on the
# Multi-IO HAT (closed = ON) gates logging, an onboard LED shows status.
ExecStart=__PYTHON__ -m watson_dms collect --switch --source hat \
  --port __PORT__ --baud __BAUD__ \
  --hat-stack __HAT_STACK__ \
  --gps-led __GPS_LED__ --logging-led __LOGGING_LED__ \
  --hat-input contact --contact-channel __CONTACT_CH__ \
  --log-dir __APPDIR__/logs --quiet
Restart=always
RestartSec=3
# Log to the journal: view with `journalctl -u watson-dms -f`.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
